# com.mendelrosenberg.screenlog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mendelrosenberg.screenlog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/mendelrosenberg/.screenlog/daemon.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mendelrosenberg/.screenlog</string>

    <key>StandardOutPath</key>
    <string>/dev/null</string>

    <key>StandardErrorPath</key>
    <string>/Users/mendelrosenberg/.screenlog/daemon.log</string>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin</string>
        <key>HOME</key>
        <string>/Users/mendelrosenberg</string>
    </dict>
</dict>
</plist>
